# shadowsxcks.conf
start on runlevel [2345]
stop on runlevel [016]
exec /{your python path}/sslocal -c /etc/{your ss config file}.js 
respawn
